# .config/systemd/user/battery-check.timer
[Unit]
Description=Run Battery Check every minute

[Timer]
OnActiveSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
